# deploy/aliyun/panelmind-frontend.service
[Unit]
Description=HireSpark Frontend
After=network.target hirespark-backend.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/hirespark/frontend
EnvironmentFile=/opt/hirespark/.env.production
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start -- --hostname 0.0.0.0 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
