# deploy/aliyun/panelmind-backend.service
[Unit]
Description=HireSpark Backend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/hirespark
EnvironmentFile=/opt/hirespark/.env.production
ExecStart=/opt/hirespark/.venv/bin/python backend/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
