# deploy/aliyun/zhiyuexingchen-frontend.service
[Unit]
Description=职跃星辰 Frontend
After=network.target zhiyuexingchen-backend.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/zhiyuexingchen/frontend
EnvironmentFile=/opt/zhiyuexingchen/.env.production
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start -- --hostname 0.0.0.0 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
